# packages/mediacenter/MuMuDVB/config/mumudvb.service
[Unit]
Description=MuMuDVB service
After=network-online.target
Before=tvlink.service

[Service]
ExecStart=/storage/.config/mumudvb/mumudvb.run 0 1 "Kyiv-t2.conf"
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=3
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
